# deploy/contract-generator.service
[Unit]
Description=Contract & Zayavka generator (Node.js)
After=network.target

[Service]
Type=simple
# Пользователь, от которого работает приложение (создайте отдельного, не root)
User=appuser
Group=appuser

# Каталог с кодом приложения
WorkingDirectory=/home/appuser/contract-generator

# Постоянное хранилище (история + справочник контрагентов).
# Лежит ВНЕ каталога кода, поэтому переживает git pull и переустановку.
Environment=DATA_DIR=/var/lib/contract-generator/data
Environment=PORT=3000
Environment=NODE_ENV=production

ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3

# Немного безопасности
NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
